# package-updates/crontabs/apt-upgrade.cron
#!/usr/bin/sh
00 6-12 * * * /home/user1/.local/bin/apt-upgrade.sh
