# package-updates/crontabs/apt-upgrade.cron
#!/usr/bin/sh
00 6-12 * * * /usr/local/bin/apt-upgrade.sh
