# package-updates/crontabs/brew-upgrade.cron
#!/bin/sh
00 6-12 * * * /Users/user1/.local/bin/brew-upgrade.sh
